# Debian/var/lib/dpkg/info/minissdpd.conffiles
/etc/default/minissdpd
/etc/init.d/minissdpd
