# Debian/var/lib/dpkg/info/minissdpd.list
/.
/etc
/etc/default
/etc/default/minissdpd
/etc/init.d
/etc/init.d/minissdpd
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/minissdpd.service
/usr
/usr/sbin
/usr/sbin/minissdpd
/usr/share
/usr/share/doc
/usr/share/doc/minissdpd
/usr/share/doc/minissdpd/README
/usr/share/doc/minissdpd/changelog.Debian.gz
/usr/share/doc/minissdpd/changelog.gz
/usr/share/doc/minissdpd/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/minissdpd.1.gz
